# config/job.cron
# Run job every sunday at 00:00

0 0 * * SUN